# process-automated-builder/scripts/systemd/process-from-flow-batch.service
[Unit]
Description=Process-from-flow batch runner (persistent)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-%h/.openclaw/.env
EnvironmentFile=%h/.config/process-from-flow-batch/env
WorkingDirectory=%h/projects/lca-skills/process-automated-builder
ExecStart=%h/projects/lca-skills/process-automated-builder/.venv/bin/python %h/projects/lca-skills/process-automated-builder/scripts/origin/process_from_flow_batch_runner.py --flow-dir ${FLOW_DIR} --state ${STATE_PATH} --log-dir ${LOG_DIR} --workers ${WORKERS} --operation ${OPERATION} --max-attempts ${MAX_ATTEMPTS} --poll-seconds ${POLL_SECONDS} --heartbeat-seconds ${HEARTBEAT_SECONDS} --stall-timeout-seconds ${STALL_TIMEOUT_SECONDS} --watch --python-bin %h/projects/lca-skills/process-automated-builder/.venv/bin/python
Restart=always
RestartSec=5
KillMode=control-group
TimeoutStopSec=30

[Install]
WantedBy=default.target
